# etc/supervisord.d/zkui_update.ini
[program:zkui-update]
command=consul-template -consul localhost:8500 -template "/etc/consul-templates/zkui.conf.ctmpl:/opt/zkui/config.cfg:supervisorctl restart zkui"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
