# etc/supervisord.d/zkui-update.ini
[program:zkui-update]
command=/opt/qnib/zkui/update/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
