# etc/supervisord.d/zkui_update.ini
[program:zkui-update]
command=/opt/qnib/zkui/bin/start_zookeeper-update.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
